# temp-throttle.service
[Unit]
Description=Throttle CPU Frequency Service

[Service]
ExecStart=/usr/sbin/temp-throttle.sh 55 57

[Install]
WantedBy=multi-user.target
